# retrovidsaverinitd

#!/bin/sh
#
# /etc/init.d/retrovidsaver
# Subsystem file for "retrovidsaver" screensaver
#
# chkconfig: 2345 95 05
# description: Screensaver that plays videos
#
# processname: retrovidsaver
# pidfile: /var/run/retrovidsaver.pid

RETVAL=0
prog="retrovidsaver"
/usr/local/bin/retrovidsaver start

start() {
    echo -n $"Starting $prog:"
    /usr/local/bin/retrovidsaver start
    RETVAL=$?
    [ "$RETVAL" = 0 ] && touch /var/lock/subsys/$prog
    echo
}

stop() {
    echo -n $"Stopping $prog:"
    /usr/local/bin/retrovidsaver stop
#    killproc $prog -TERM
    RETVAL=$?
    [ "$RETVAL" = 0 ] && rm -f /var/lock/subsys/$prog
    echo
}

reload() {
    echo -n $"Reloading $prog:"
    killproc $prog -HUP
    RETVAL=$?
    echo
}

case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    restart)
        stop
        start
        ;;
    reload)
        reload
        ;;
    condrestart)
        if [ -f /var/lock/subsys/$prog ] ; then
            stop
            # avoid race
            sleep 3
            start
        fi
        ;;
    status)
        status $prog
        RETVAL=$?
        ;;
    *)
        echo $"Usage: $0 {start|stop|restart|reload|condrestart|status}"
        RETVAL=1
esac
exit $RETVAL
